# supervisord/conf.d/gluetun.conf
[program:gluetun]
command         = docker compose -f /opt/hermes/stacks/vpn/docker-compose.yml --env-file /opt/hermes/stacks/vpn/.env up --no-log-prefix --no-color gluetun
directory       = /opt/hermes/stacks/vpn
priority        = 10
autostart       = true
autorestart     = true
startretries    = 5
startsecs       = 15
stopwaitsecs    = 15
stdout_logfile  = /var/log/supervisor/gluetun.stdout.log
stderr_logfile  = /var/log/supervisor/gluetun.stderr.log
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile_backups  = 3
stderr_logfile_backups  = 3
